# private_dot_config/autostart/ssh-proxy-tray.desktop
[Desktop Entry]
Name=SSH Proxy Tray
Comment=Monitor and control SSH SOCKS5 proxy connection
Exec=/usr/bin/python3 /path/to/ssh_proxy_tray.py
Icon=network-vpn
Terminal=false
Type=Application
Categories=Network;Utility;
